# res/files/98-backup.rules
KERNEL=="sd?1", ACTION=="add", SUBSYSTEMS=="usb", RUN+="/home/rbe/GN01-Update.sh %k"
